# load-from-static.cron
0 */6 * * * docker compose --file /srv/docker/life-ustc/docker-compose.yml exec --user root -it app bun run tools/load-from-static.ts